# supervisor/redis.conf
[program:redis]
directory=/opt/redis/scripts
command=python3 -m run_redis
autorestart=true
startretries=10
startsecs=6
