# supervisor/require-master.conf
[program:require-master]
directory=/opt/redis/scripts
command=python -m require_master
startsecs=0
