# supervisor/run-redis.conf
[program:redis]
command=node /opt/redis/app.js
autorestart=true
startretries=10
startsecs=6